# root/etc/services.d/jdownloader/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc \
	java -Dsun.jnu.encoding=UTF-8 -Dfile.encoding=UTF-8 -Djava.awt.headless=true -jar /opt/JDownloader/JDownloader.jar -norestart
